# archlive/work/x86_64/airootfs/var/lib/pacman/local/virtualbox-guest-utils-5.0.10-1/files
%FILES%
etc/
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/vboxclient.desktop
usr/
usr/bin/
usr/bin/VBoxClient
usr/bin/VBoxClient-all
usr/bin/VBoxControl
usr/bin/VBoxService
usr/bin/mount.vboxsf
usr/lib/
usr/lib/VBoxOGL.so
usr/lib/VBoxOGLarrayspu.so
usr/lib/VBoxOGLcrutil.so
usr/lib/VBoxOGLerrorspu.so
usr/lib/VBoxOGLfeedbackspu.so
usr/lib/VBoxOGLpackspu.so
usr/lib/VBoxOGLpassthroughspu.so
usr/lib/security/
usr/lib/security/pam_vbox.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/vboxservice.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/virtualbox-guest-utils.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-vboxguest.rules
usr/lib/xorg/
usr/lib/xorg/modules/
usr/lib/xorg/modules/dri/
usr/lib/xorg/modules/dri/vboxvideo_dri.so
usr/lib/xorg/modules/drivers/
usr/lib/xorg/modules/drivers/vboxvideo.so
usr/share/
usr/share/licenses/
usr/share/licenses/virtualbox-guest-utils/
usr/share/licenses/virtualbox-guest-utils/LICENSE

